# .notes/frogcam.conf
[program:frogcam-flask]
directory=/home/pi/frogcam/app
environment=PYTHONPATH=/usr/lib/python3/dist-packages
command=/home/pi/frogcam/app/.venv/bin/gunicorn --bind 127.0.0.1:8000 wsgi:app
autostart=true
autorestart=true
stderr_logfile=/var/log/frogcam/wsgi.err.log
stdout_logfile=/var/log/frogcam/wsgi.out.log

[program:frogcam-ws]
directory=/home/pi/frogcam/app
environment=PYTHONPATH=/usr/lib/python3/dist-packages
command=/home/pi/frogcam/app/.venv/bin/python ws.py --bind 127.0.0.1:8001
autostart=true
autorestart=true
stderr_logfile=/var/log/frogcam/ws.err.log
stdout_logfile=/var/log/frogcam/ws.out.log
